# Ansible_Tower/Playbook_Examples/SECURITY/ipa-audit-master/files/audit_audit.timer
[Unit]
Description=Audit Shipping Timer

[Timer]
OnBootSec=15min
OnCalendar=*:0/15
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
