# Ansible_Tower/Playbook_Examples/LAB/lab-system-monitoring-master/external-configs/etc/cron.hourly/health-check
#!/bin/sh

/usr/sbin/logrotate --force /etc/logrotate.d/health-check &> /dev/null
