# Ansible_Tower/Playbook_Examples/LAB/ansible-lab-master/prorail/files/ip-add.service
[Unit]
Description=Add extra IPs to host
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ip-add.sh

[Install]
WantedBy=default.target
